# packaging/macos/com.windsentinel.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.windsentinel.agent</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Library/WindSentinel/bin/windsentinel_agent</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>WINDSENTINEL_AGENT_CONFIG_PATH</key>
      <string>/Library/WindSentinel/config/config.json</string>
      <key>WINDSENTINEL_AGENT_STATE_DIR</key>
      <string>/Library/WindSentinel/state</string>
      <key>WINDSENTINEL_CONFIG_VERIFY_KEY_B64</key>
      <string>__VERIFY_KEY_B64__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Library/WindSentinel/logs/agent.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/WindSentinel/logs/agent.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Library/WindSentinel</string>
  </dict>
</plist>
